# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/tmp/supervisord.pid

; ─── nginx ────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ─── Express API server ───────────────────────────────────────────
[program:api]
command=node --enable-source-maps /app/artifacts/api-server/dist/index.mjs
directory=/app
autostart=true
autorestart=true
; Delay API start slightly so DB connection attempts don't race with startup
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
; Environment variables are inherited from the container env.
; Make sure these are set in your HF Space secrets:
;   PORT, DATABASE_URL, CLERK_SECRET_KEY, GOOGLE_APPLICATION_CREDENTIALS_JSON
environment=PORT="8080",NODE_ENV="production"
